# trunk/src/initrd/contents/etc/inittab
# Формат файла:
# программа:аргументы:текст для вывода
# если программа не принимает аргументы или не требуется выводить текст - поля можно не заполнять,
# но нужно сохранить все разделители.
#
/bin/i8042::Starting keyboard/mouse driver
#/bin/pci_scan::Scanning PCI bus
/usr/bin/fgs::Starting up GUI
#/bin/shell::Starting up shell
