# trunk/src/initrd/contents/etc/inittab
/bin/stty::
.:open_tty:/dev/tty
/bin/i8042::Starting keyboard/mouse driver
/bin/pcid::Starting PCI daemon
/bin/pci_scan::Scanning PCI bus
/usr/bin/gwinsy::Starting up GUI
/bin/rtc::Starting RTC
/bin/tmpfs::Starting up tmpfs
