# trunk/src/initrd/contents/etc/RLY/inittab
/bin/Keyboard;service;Starting up keyboard driver
/bin/cat;/etc/ascii;[silent]
/bin/gsh;-l;Starting shell