# src/initrd/contents/etc/inittab
/bin/i8042::Starting keyboard/mouse driver
#/bin/pci_scan::Scanning PCI bus
/usr/bin/fgs::Starting up GUI
#/bin/shell::Starting up shell
/bin/rtc::Starting RTC
/bin/tmpfs::Starting up tmpfs
